# vagrant-setup/crontab
* * * * * /vagrant/app/console app:clear_old_clients_personal_info_command
